# cvm/cert-manager/supervisord.conf
# TODO: we are using a custom logger to differentiate logs from different processes
# This will be solved by https://github.com/Supervisor/supervisor/pull/1689
# Then we can remove stdout_logfile_maxbytes=0 and the proc-logger program

[supervisord]
nodaemon=true # start in foreground
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=0
priority=100 # lower start first

[program:cert-manager]
directory=/app
command=uv run cert_manager.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cert-manager.log
stdout_logfile_maxbytes=0
priority=200

[program:proc-logger]
directory=/app
command=uv run proc_logger.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=300

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
