# cvm/cert-manager/supervisord.conf
[supervisord]
nodaemon=true # start in foreground
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx_stderr.log
stdout_logfile=/var/log/supervisor/nginx_stdout.log
priority=100 # lower start first

[program:cert-manager]
directory=/app
command=uv run cert_manager.py
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cert-manager_stderr.log
stdout_logfile=/var/log/supervisor/cert-manager_stdout.log
priority=200

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
